# deploy/piv-agent.service
[Unit]
Description=piv-agent service

[Service]
ExecStart=piv-agent serve --debug --agent-types=ssh=0;gpg=1
