# deploy/piv-agent.service
[Unit]
Description=piv-agent service

[Service]
ExecStart=%h/go/bin/piv-agent serve --debug
