# deploy/piv-agent.service
[Unit]
Description=piv-agent service

[Service]
ExecStart=%h/dev/piv-agent/piv-agent serve --debug
