# deploy/piv-agent.socket
[Unit]
Description=piv-agent socket

[Socket]
ListenStream=%t/piv-agent/ssh.socket

[Install]
WantedBy=sockets.target
